# config/autostart/utzmachine.desktop
[Desktop Entry]
Type=Application
Exec=/usr/bin/utz
